# supervisord-ttyd-sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D

default_startsecs=0

default_autostart=true

default_autorestart=true

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1024x768x16
autostart=true
autorestart=true

[program:x11vnc]
command=x11vnc -forever -usepw -rfbauth /etc/x11vnc.pass -display :0
autostart=true
autorestart=true

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 8080 localhost:5900
autostart=true
autorestart=true 