# supervisord-ttyd-sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D

[program:ttyd]
command=ttyd --writable -p 7681 bash 